# src/systemd/docker-filevol-plugin.service
[Unit]
Description=Docker Volume driver for file volume storage
Before=docker.service

[Service]
Type=simple
EnvironmentFile=-/etc/docker/filevol-plugin
ExecStart=/usr/libexec/docker/docker-filevol-plugin
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
